# discord_bot.service
/*

The Discord bot will run as a persistent service on the EC2 instance. The Terraform script includes a systemd service file to ensure the bot starts automatically upon system boot and remains running.

This service file is automatically deployed through Terraform to the EC2 instance running the Discord bot.
The bot will start on boot and run continuously in the background.

*/


[Unit]
Description=Discord Bot Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 /path/to/discord_bot.py
WorkingDirectory=/path/to/bot-directory
Restart=always
User=ec2-user

[Install]
WantedBy=multi-user.target
